# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid

; ---------------------------------------------------------------
; nginx — serves the React static build and proxies /api/ traffic
; ---------------------------------------------------------------
[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=unexpected
exitcodes=0
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ---------------------------------------------------------------
; api — gunicorn + UvicornWorker (multi-process ASGI)
; GUNICORN_WORKERS is injected from the container environment.
; --forwarded-allow-ips=* trusts the X-Forwarded-* headers that
; nginx sets, so FastAPI sees the real client IP and protocol.
; ---------------------------------------------------------------
[program:api]
command=/venv/bin/gunicorn src.main:app
    -k uvicorn.workers.UvicornWorker
    -w %(ENV_GUNICORN_WORKERS)s
    -b 127.0.0.1:8888
    --forwarded-allow-ips=*
    --access-logfile -
    --error-logfile -
directory=/app/backend
autostart=true
autorestart=unexpected
exitcodes=0
stopasgroup=true
killasgroup=true
priority=20
user=www-data
environment=PYTHONUNBUFFERED="1",PYTHONDONTWRITEBYTECODE="1",HOME="/var/www"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
